# ansible/templates/supervisord.conf
[supervisord]
nodaemon=true

[program:twitter-worker]
command={{ symfony_console_path }} process:queue --cache=file_system --env=prod
autostart=true
autorestart=true